# dotf/crontab
0 10 * * 1 cd ziprecruiter && git fetch --prune --prune-tags
